# conf/filter.ini
[program:filter]
command=python /home/liuweibo/skyline/skyline-client/filter.py /home/liuweibo/skyline/conf//monitor.conf /home/liuweibo/skyline/conf//warning.conf
autostart=false                ; start at supervisord start (default: true)
