# conf/gather.ini
[program:gather]
command=python /home/liuweibo/skyline/skyline-client/gather.py /home/liuweibo/skyline/conf//monitor.conf
autostart=false                ; start at supervisord start (default: true)
